# deploy/homebase-nas-agent.service
[Unit]
Description=Homebase NAS Host Agent
After=network.target

[Service]
Type=simple
User=homebase
Group=homebase
WorkingDirectory=/opt/homebase-nas-agent
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=-/etc/homebase/nas-agent.env

[Install]
WantedBy=multi-user.target
